# scripts/deploy/pejortryia.service
[Unit]
Description=Pejortryia Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bots
Group=bots
WorkingDirectory=/opt/pejortryia
EnvironmentFile=/opt/pejortryia/.env
ExecStart=/usr/bin/uv run python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
